# systemd-units/smartctl-long@.service
[Unit]
Description=Smartctl long test of %i

[Service]
Type=oneshot
ExecStart=/usr/bin/smartctl -t long %i
