# systemd-units/systemd-poweroff@.timer
[Unit]
Description=Poweroff the system at the given time

[Timer]
OnCalendar=%I
Unit=systemd-poweroff.service

[Install]
WantedBy=timers.target
